# packaging/linux/systemd/cloudtolocalllm-bridge.service
[Unit]
Description=CloudToLocalLLM Desktop Bridge
Documentation=https://cloudtolocalllm.online
After=network.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=%i
ExecStart=/opt/cloudtolocalllm-bridge/cloudtolocalllm-bridge --daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/%i/.config/cloudtolocalllm

# Environment
Environment=HOME=/home/%i
Environment=XDG_CONFIG_HOME=/home/%i/.config

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudtolocalllm-bridge

[Install]
WantedBy=default.target
