# config/systemd/cloudtolocalllm-tray.service
[Unit]
Description=CloudToLocalLLM System Tray Daemon
Documentation=https://cloudtolocalllm.online
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=%i
Group=%i
ExecStart=/usr/bin/cloudtolocalllm-tray
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/%i/.cloudtolocalllm

# Environment
Environment=HOME=/home/%i
Environment=XDG_CONFIG_HOME=/home/%i/.config
Environment=DISPLAY=:0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudtolocalllm-tray

[Install]
WantedBy=default.target
