# scripts/cloudtolocalllm.service
[Unit]
Description=CloudToLocalLLM Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/cloudtolocalllm/portal
ExecStartPre=/bin/bash -c 'mkdir -p /opt/cloudtolocalllm/logs'

# Ensure Docker is running
ExecStartPre=/bin/bash -c 'systemctl is-active --quiet docker || systemctl start docker'

# Start services
ExecStart=/usr/bin/docker-compose -f docker-compose.auth.yml -f docker-compose.web.yml up

# Stop services
ExecStop=/usr/bin/docker-compose -f docker-compose.auth.yml -f docker-compose.web.yml down

# Restart policy
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 